# crontab
0 7 * * * python3 /src/maxmeter-test.py
