# crontab
0 6 * * * 1000720000 python3 /src/maxmeter-test.py > /proc/1/fd/1 2>/proc/1/fd/2
